# etc/services.d/duplicati/run
#!/usr/bin/execlineb -P
with-contenv importas PORT PORT
/usr/bin/duplicati-server --webservice-port=${PORT} --webservice-interface=any
